# src/systemd-files/GOLDiProgrammingService.socket
[Unit]
Description = Socket for the GOLDi Programming Service

[Socket]
ListenStream = @GOLDiProgrammingService

[Install]
WantedBy = sockets.target
