# src/systemd-files/GOLDiCommandService.service
[Unit]
Description = The GOLDi Command Service

[Service]
ExecStart = GOLDiCommandService

[Install]
WantedBy = multi-user.target
