# src/systemd-files/GOLDiWebcamService.service
[Unit]
Description = The GOLDi Webcam Service

[Service]
ExecStart = GOLDiWebcamService

[Install]
WantedBy = multi-user.target
