# src/systemd-files/GOLDiCommandService.socket
[Unit]
Description = Socket for the GOLDi Command Service

[Socket]
ListenStream = @GOLDiCommandService

[Install]
WantedBy = sockets.target
