# src/systemd-files/GOLDiProgrammingService.service
[Unit]
Description = The GOLDi Programming Service

[Service]
ExecStart = GOLDiProgrammingService

[Install]
WantedBy = multi-user.target
