# src/systemd-files/GOLDiProtectionService.service
[Unit]
Description = The GOLDi Protection Service

[Service]
ExecStart = GOLDiProtectionService

[Install]
WantedBy = multi-user.target
