# src/systemd-files/GOLDiCommunicationService.service
[Unit]
Description = The GOLDi Communication Service

[Service]
ExecStart = GOLDiCommunicationService

[Install]
WantedBy = multi-user.target
