# src/systemd-files/GOLDiInitializationService.service
[Unit]
Description = The GOLDi Initialization Service

[Service]
ExecStart = GOLDiInitializationService

[Install]
WantedBy = multi-user.target
