# src/systemd-files/GOLDiProtectionService.socket
[Unit]
Description = Socket for the GOLDi Protection Service

[Socket]
ListenStream = @GOLDiProtectionService

[Install]
WantedBy = sockets.target
