# src/systemd-files/GOLDiWebcamService.socket
[Unit]
Description = Socket for the GOLDi Webcam Service

[Socket]
ListenStream = @GOLDiWebcamService

[Install]
WantedBy = sockets.target
